# files/etc/uci-defaults/99-custom
uci -q set wireless.@wifi-device[0].disabled='1'
uci -q set wireless.@wifi-iface[0].ssid='OpenWrt0815'